# scripts/mettaton.service
[Unit]
Description=Mettaton Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/mettaton
ExecStart=/usr/bin/node src/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mettaton

# Environment
EnvironmentFile=/home/pi/mettaton/.env

# Resource limits for Raspberry Pi 3B+
MemoryMax=512M

[Install]
WantedBy=multi-user.target
